# salt-master/supervisor-salt.conf
[program:salt-master]
command=/usr/bin/salt-master
numprocs=1
autostart=true
autorestart=true

[program:sshd]
command=/usr/sbin/sshd -D
numprocs=1
autostart=true
autorestart=true

