# salt-minion/supervisor-salt.conf
[program:salt-minion]
command=/usr/bin/salt-minion
numprocs=1
autostart=true
autorestart=true

